# packages/sdr-host/docker/overlay/s6-overlay/s6-rc.d/sdr-host-api/run
#!/bin/sh
# Wait for rtlmux to be ready
sleep 1

# Start the Node.js API server
cd /app
exec node dist/index.js
