# docker/overlay/s6-overlay/s6-rc.d/wavekit-init/up
#!/bin/sh
exec /usr/local/bin/init-system.sh
