# docker/overlay/s6-overlay/s6-rc.d/wavekit-api/finish
#!/bin/sh
/usr/local/bin/finish-api.sh "$@"
exit 0
