# docker/overlay/s6-overlay/s6-rc.d/wavekit-api/finish
#!/bin/sh
# Graceful shutdown for WaveKit API

exec chpst -u root:root s6-svwait -t 10000 -d /run/service/wavekit-api 2>/dev/null || true
