# packages/sdr-host/docker/overlay/s6-overlay/s6-rc.d/sdr-host-base/up
#!/bin/sh
# Base init - create required directories
mkdir -p /var/log/sdr-host /var/run/sdr-host
