# packages/sdr-host/docker/overlay/s6-overlay/s6-rc.d/rtlmux/run
#!/command/with-contenv sh
# Wait for rtl_tcp to be ready
sleep 2

IQ_PORT=${SDR_HOST_RTLMUX__PORT:-5555}

# Start rtlmux (stats served on IQ port + 1)
exec rtlmux \
    -h 127.0.0.1 \
    -p ${SDR_HOST_RTL_TCP__INTERNAL_PORT:-1234} \
    -l ${IQ_PORT}
