# docker/overlay/s6-overlay-sdrpp/s6-rc.d/sdrpp-server/finish
#!/bin/sh
/usr/local/bin/finish-sdrpp.sh "$@"
exit 0
