# docker/overlay/s6-overlay/s6-rc.d/sdrpp-server/run
#!/bin/sh
# s6-rc service: sdrpp-server
# Purpose: SDR++ server for IQ stream capture from hardware
# Upstream: RTL-TCP on Pi or direct USB device

exec 2>&1
exec chpst -u root:root \
  sdrpp \
  --server \
  --port 5259 \
  --addr 0.0.0.0 \
  --log-level info
